# amdgpuoc.service
[Unit]
Description = Overclock an AMD GPU

[Service]
ExecStart = /usr/local/sbin/amdgpuoc
Type = oneshot

[Install]
WantedBy = multi-user.target
