# amdgpuoc.service
[Unit]
Description = Overclock an AMD GPU

[Service]
ExecStart = /usr/sbin/amdgpuoc

[Install]
WantedBy = multi-user.target
